# build_files/files/common/usr/lib/systemd/system/install-default-flatpaks.service
[Unit]
Description=Install default Flatpak apps on first boot
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/install-default-flatpaks/done

[Service]
Type=oneshot
RemainAfterExit=yes
StateDirectory=install-default-flatpaks
ExecStart=/usr/libexec/install-default-flatpaks
ExecStartPost=/usr/bin/touch /var/lib/install-default-flatpaks/done

[Install]
WantedBy=multi-user.target
